# chapter-06/examples/upstart/nodebook.conf
start on started network-services

chdir /usr/local/node-app

exec /usr/bin/npm start
