# chapter-06/examples/upstart/nodebook.conf
start on started network-services

chdir /usr/local/node-app

respawn
respawn limit 5 120

exec /usr/bin/npm start
